# templates/kong.init.erb
#!/bin/sh
#
# chkconfig: - 80 45
# description: Starts and stops Kong
# update deamon path to point to the kong executable

start() {
        echo -n "Starting Kong"
        <%= @kong_install_path + "kong"%> start -c /etc/kong/kong.yml
}

stop() {
        echo -n "Stopping Kong"
        <%= @kong_install_path + "kong"%> stop -c /etc/kong/kong.yml
}

restart() {
        echo -n "Restarting Kong"
        <%= @kong_install_path + "kong"%> restart -c /etc/kong/kong.yml
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 1
esac

exit $?